# nfs-qbit.service
[Unit]
Description=NFS Mount and qBittorrent Startup Service
After=network.target docker.service
Wants=docker.service

[Service]
Type=oneshot
ExecStart=/root/start-nfs-and-qbit.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal
User=root
Group=root

# 重启策略
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target